# devops/docker/files/etc/services.d/java/run
#!/usr/bin/with-contenv sh
/opt/openjdk-15/bin/java -jar /app.jar --spring.config.additional-location=/application.properties
